# config/supervisor-cals.conf
[program:cals]
command=python3 cals.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
