# config/supervisor-news.conf
[program:news]
command=python3 news.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
