# config/supervisor-cals-detail.conf
[program:cals_detail]
command=python3 cals_detail.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
