# config/supervisor-res.conf
[program:res]
command=python3 res.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
