# config/supervisor-rates.conf
[program:rates]
command=python3 rates.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
