# config/supervisor-arts.conf
[program:arts]
command=python3 article.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
