# config/supervisor.conf
[program:fx]
command=/usr/local/bin/gunicorn wsgi -c config/gunicorn.py
directory=/var/www/fx-calendar
autostart=true
autorestart=true
